# bbs.conf
[program:bbs]
command=/dev/pydev/bbs/bin/gunicorn wsgi -c gunicorn.config.py
directory=/var/www/bbs
autostart=false
autorestart=false